# systemd/unicorn.service.erb
[Unit]
Description=<%= description %>
After=<%= after.join(' ') %>
Wants=<%= wants.join(' ') %>

[Service]
Type=forking
User=<%= user %>
Environment=BUNDLE_GEMFILE=<%= gemfile_path %>
WorkingDirectory=<%= working_directory %>
ExecStart=/usr/bin/bundle exec unicorn -c <%= config_path %> -E production -D
ExecReload=/bin/kill -HUP $MAINPID
ExecRestart=/bin/kill -USR2 $MAINPID
ExecRestart=/usr/bin/sleep <%= restart_wait_seconds %>
ExecRestart=/bin/kill -QUIT $MAINPID
ExecStop=/bin/kill -QUIT $MAINPID
Restart=always
PIDFile=<%= pid_path %>

[Install]
WantedBy=multi-user.target
